# crontab/port_check.sh
#!/bin/bash


curl http://localhost:88/port_check

